# root/etc/systemd/system/s6-overlay.service
[Unit]
Description=S6 overlay


[Service]
ExecStart=/init

[Install]
WantedBy=multi-user.target